# supervisord.conf
# =============================================================================
# NKDash Solo Mode — supervisord configuration
# =============================================================================
# Manages 3 processes in a single container:
#   1. gunicorn   — Dash BI dashboard (public, port 8050)
#   2. streamlit   — Admin UI (internal, port 8501)
#   3. scheduler   — ETL daemon (background, no port)
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null           ; log to stdout (docker logs captures it)
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[program:gunicorn]
command=gunicorn
    --workers 1
    --threads 4
    --timeout 180
    --graceful-timeout 30
    --bind 0.0.0.0:8050
    --access-logfile -
    --error-logfile -
    app:server
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:streamlit]
command=streamlit run admin/app.py
    --server.port 8501
    --server.address 0.0.0.0
    --server.headless true
    --server.enableCORS false
    --server.enableXsrfProtection false
    --browser.gatherUsageStats false
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:scheduler]
command=python -u scheduler/main.py
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30
